# deps/beacon.service
[Unit]
Description=beacon

[Service]
ExecStart=/usr/bin/rkt run \
	--net=host \
        --volume redis,kind=host,source=/tmp/dump.rdb \
        unixvoid.com/beacon

[Install]
WantedBy=multi-user.target
